# infra/pi-image/systemd/discodb2-backend.service
[Unit]
Description=discodb2 — thin CAN backend (stream/record/replay, listen-only)
Documentation=file:/opt/discodb2/docs/DESIGN.md
# Start after the bus and the AP so a phone that connects on boot sees data.
# These are Wants (not Requires): the backend must still run with sim/replay
# even if can0 is unplugged — DESIGN §4.4 (sim works with zero hardware).
Wants=discodb2-can.service discodb2-ap.service
After=discodb2-can.service discodb2-ap.service network.target

[Service]
Type=simple
EnvironmentFile=/etc/discodb2/discodb2.env
# The package is `discodb2_backend`, living at backend/discodb2_backend/, so the
# import root (and working dir) is the backend/ dir. This mirrors the working
# docker entrypoint (`cd /app/backend && python -m discodb2_backend ...`).
WorkingDirectory=/opt/discodb2/backend
# Launch the backend package and pass ONLY the flags config.py actually defines
# (--source/--bitrate/--port/--host). There is NO --listen-only flag: listen-only
# is enforced unconditionally server-side and at the link layer by can-up.sh.
ExecStart=${BACKEND_PYTHON} -m discodb2_backend --source ${BACKEND_SOURCE} --bitrate ${CAN_BITRATE} --port ${BACKEND_PORT} --host 0.0.0.0
Restart=on-failure
RestartSec=2
# Run unprivileged where possible; SocketCAN read needs CAP_NET_RAW only.
User=discodb2
AmbientCapabilities=CAP_NET_RAW
NoNewPrivileges=yes
# Hardening that does not interfere with SocketCAN/USB read.
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/discodb2/recordings
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
